# saas_bot/deploy/saas-bot-api.service
[Unit]
Description=SaaS Telegram Bot Admin API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/saas_bot
EnvironmentFile=/opt/saas_bot/.env
ExecStart=/opt/saas_bot/.venv/bin/uvicorn api.server:app --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/saas_bot/api.log
StandardError=append:/var/log/saas_bot/api.err.log

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/saas_bot/data /var/log/saas_bot

[Install]
WantedBy=multi-user.target
